# java-app/crontab
0 2 * * * cd /app && java -jar /app/app.jar >> /proc/1/fd/1 2>> /proc/1/fd/2 && touch /app/output/report_$(date +\%d-\%m-\%Y).pdf 